# root/etc/local.d/grayscale.start
#!/bin/sh
nvidia-settings -a DigitalVibrance=-1024 > /dev/null
